# .config/systemd/user/mate-power-manager.service
[Unit]
Requires=xsession.target

[Service]
Type=simple
ExecStart=/usr/bin/mate-power-manager
Restart=on-failure

[Install]
WantedBy=xsession.target
